# examples/midori/rootfs/etc/services.d/midori/run
#!/usr/bin/with-contenv sh

echo "> Starting Midori Web Browser"

exec /usr/bin/midori "https://github.com/N0rthernL1ghts/docker-app"