# examples/midori/rootfs/etc/services.d/midori/run
#!/usr/bin/with-contenv sh

echo "> Starting Midori Web Browser"

s6-setuidgid app /usr/bin/midori -a "https://github.com/N0rthernL1ghts/docker-app"